# roots/templates/ceph/ceph-mon.service
[Unit]
Description=ceph monitor process

[Service]
EnvironmentFile=/etc/ceph/envfile
ExecStart=/usr/bin/ceph-mon -f -c ${CONF} -i ${FQDN}
