# roots/templates/ceph/osd-mountpoint.mount
[Mount]
What={{ part }}
Where={{ path }}
Type=ext4
Options=user_xattr

[Install]
WantedBy=multi-user.target
